# ansible/roles/borg_backup_users/files/borg-backup.service
[Unit]
Description=Back up files with Borg

[Service]
Type=oneshot
ExecStart=%h/.local/bin/borg-backup.sh
